# crontab
7 0,12 * * * /usr/local/bin/node /index.js
